# templates/centos8-container/supervisord/supervisord.d/opennebula-ssh-agent.ini
[program:opennebula-ssh-agent]
command=/usr/share/one/supervisord/service/opennebula-ssh-agent.sh
user=oneadmin
directory=/var/lib/one
priority=201
autostart=false
autorestart=true
startsecs=5
exitcodes=0,2
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0

